# deploy/canquery-worker.service
# /etc/systemd/system/canquery-worker.service
# Long-running ingest worker (systemd rather than tmux, per the runbook).
[Unit]
Description=canquery ingest worker (drains ingest_jobs)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=canquery
Group=canquery
WorkingDirectory=/home/canquery/canquery/server
ExecStart=/usr/bin/node scripts/ingest-worker.js
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=10

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/canquery/canquery/server /tmp
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
