# deploy/opencanada-worker.service
# /etc/systemd/system/opencanada-worker.service
# Long-running ingest worker (systemd rather than tmux, per the runbook).
[Unit]
Description=opencanada ingest worker (drains ingest_jobs)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=opencanada
Group=opencanada
WorkingDirectory=/home/opencanada/opencanada/server
ExecStart=/usr/bin/node scripts/ingest-worker.js
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=10

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/opencanada/opencanada/server /tmp
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
